# etc/ss/ss-init
#!/bin/sh /etc/rc.common
# Copyright hiwifi.com

#should start before dnsmasq = 47
START=46

appname=SSAPPNAME
ss_lib=/etc/ss/ss-lib

source $ss_lib

start() {
	ss_start
}

stop() {
	ss_stop
}

