# usr/lib/opkg/info/geewan-ss.list
/etc/gw-shadowsocks/gw-shadowsocks.dnslist
/etc/init.d/gw-redsocks
/etc/init.d/gw-shadowsocks
/etc/rc.d/S99gw-shadowsocks
/lib/gw-shadowsocks.sh
/lib/ss-test.lua
/usr/bin/dns2socks
/usr/bin/ss-local
/usr/bin/ss-redir
/usr/bin/ss-tunnel
/usr/lib/lua/geewan.lua
/usr/lib/lua/luci/log.lua
/usr/lib/lua/luci/controller/prometheus/shadowsocks.lua
/usr/lib/lua/luci/controller/api/prometheus.lua
/usr/lib/lua/luci/view/admin_web/prometheus/index.htm
/usr/lib/lua/luci/view/admin_web/menu/menu_left.htm
