# usr/lib/opkg/info/geewan-ss.list
/etc/gw-shadowsocks/gw-shadowsocks.dnslist
/etc/init.d/gw-redsocks
/etc/init.d/gw-shadowsocks
/etc/rc.d/S80gw-shadowsocks
/lib/gw-shadowsocks.sh
/lib/ss-test.lua
/usr/bin/dns2socks
/usr/bin/ss-local
/usr/bin/ss-redir
/usr/bin/ss-tunnel
/usr/lib/lua/geewan.lua
/usr/lib/lua/luci/controller/admin_web/shadowsocks.lua
/usr/lib/lua/luci/controller/api/geewan.lua
/usr/lib/lua/luci/view/admin_web/menu/adv_menu.htm
/usr/lib/lua/luci/view/admin_web/network/cfg_import.htm
/usr/lib/lua/luci/view/admin_web/network/shadowsocks.htm
/usr/lib/lua/luci/view/admin_web/network/ss_adv.htm
