# usr/lib/opkg/info/geewan-ss.list
/etc/init.d/gw-redsocks
/etc/init.d/gw-shadowsocks
/etc/rc.d/S99gw-shadowsocks
/etc/gw-shadowsocks/base-gfwlist.txt
/etc/gw-shadowsocks/china-banned
/etc/gw-shadowsocks/upchinaroute.sh
/etc/gw-shadowsocks/upgfwlist.sh
/lib/gw-shadowsocks.sh
/lib/ss-test.lua
/usr/bin/dns2socks
/usr/bin/ss-local
/usr/bin/ss-redir
/usr/bin/ss-tunnel
/usr/bin/ssr-local
/usr/bin/ssr-redir
/usr/bin/ssr-tunnel
usr/bin/base64
/usr/lib/libpcre.so.1
/usr/lib/libsodium.so.18
/usr/lib/lua/luci/controller/admin_web/shadowsocks.lua
/usr/lib/lua/luci/controller/api/prometheus.lua
/usr/lib/lua/luci/controller/prometheus/shadowsocks.lua
/usr/lib/lua/luci/view/admin_web/menu/menu_left.htm
/usr/lib/lua/luci/view/admin_web/prometheus/index.htm
/usr/lib/lua/luci/log.lua